# root/etc/services.d/mount/run
#!/usr/bin/with-contenv sh
. "/usr/bin/variables"

umask 022

mount_command="/usr/bin/mount"

echo "Started ${mount_command}"
exec $mount_command
